# lib/systemd/system/beldex-storage-servers.target
[Unit]
Description=beldex storage servers
